# gentility-agent.service
[Unit]
Description=Gentility AI Agent Daemon
After=network.target

[Service]
ExecStart=/usr/bin/gentility-agent
Restart=always
RestartSec=5
User=gentility
WorkingDirectory=/var/lib/gentility-agent
Environment=PATH=/usr/local/bin:/usr/bin:/bin
EnvironmentFile=-/etc/gentility.conf

[Install]
WantedBy=multi-user.target